# sondehub/rootfs/etc/services.d/sondehub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the SondeHub add-on service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
exec python3 /usr/bin/sondehub_addon.py
